# scripts/weather_daemon/service/weather-daemon.timer
[Unit]
Description=Run VATSIM Weather Briefing Generator every 5 minutes at XX:54/59/04/09...
Documentation=https://github.com/leftos/vatsim_control_recs

[Timer]
# Run at fixed times: XX:04, XX:09, XX:14, ... XX:54, XX:59
# This aligns with METAR updates which typically occur at XX:53
OnCalendar=*:4/5

# If system was off, don't catch up on missed runs
Persistent=false

# Accuracy allows up to 1 second of delay for batching
AccuracySec=1s

[Install]
WantedBy=timers.target
